# backend/supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/php-fpm.log
